# agent.service
[Unit]
Description=System Monitoring Agent
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/agent.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target
